# systemd/ads-optimizer-light.service
[Unit]
Description=OwnersHub Ad & Page Optimiser - midweek light pass
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=rodean
WorkingDirectory=/home/rodean/ads-optimizer
ExecStart=/home/rodean/ads-optimizer/run.sh light
RuntimeMaxSec=600
Nice=10
# OnFailure=ads-optimizer-failure-email.service

[Install]
WantedBy=multi-user.target
